# init_script/nas_control.service
[Unit]
Description=Nas control

[Service]
User=nico
Type=simple
ExecStart=/home/nico/Nas_Control/nas_control.py

[Install]
WantedBy=multi-user.target